# bamf/root/etc/services.d/huey/run
#!/usr/bin/with-contenv bash

exec /virt/bin/python3 /virt/manage.py run_huey
