# bamf/root/etc/services.d/django/run
#!/usr/bin/with-contenv bash

exec /virt/bin/python3 /virt/manage.py runserver 0.0.0.0:8000
