# bamf/root/etc/services.d/redis/run
#!/usr/bin/with-contenv bash

exec /usr/local/bin/redis-server /virt/redis.conf
